# services/robot.service
[Unit]
Description=robot service

[Service]
User=ubuntu
ExecStart=/bin/bash --login -c 'source /home/ubuntu/.bashrc; /home/ubuntu/mini_pupper_ros_bsp/services/robot.sh'
KillSignal=2
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
